# docker-compose/rocketchat/rocketchat.service
[Unit]
Description=rocketchat
Requires=docker.service
Requires=mongo.service
After=docker.service
After=mongo.service

[Service]
EnvironmentFile=/etc/environment
EnvironmentFile=/etc/default/rocketchat
User=root
Type=oneshot
RemainAfterExit=true
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker pull rocketchat/rocket.chat:${ROCKETCHAT_VERSION}

ExecStart=/usr/local/bin/docker-compose -f /opt/admin-tools/docker-compose/rocketchat/docker-compose.yml --env-file /etc/default/rocketchat up -d --remove-orphans
ExecStop=/usr/local/bin/docker-compose -f /opt/admin-tools/docker-compose/rocketchat/docker-compose.yml --env-file /etc/default/rocketchat down

[Install]
WantedBy=multi-user.target

