# munin-asyncd/munin-async.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=/opt/admin-tools/munin-asyncd/munin-asyncd --verbose --timeout 100 --interval 3600 --retain 48 --rate 100